# deploy/supervisord-uwsgi.ini
[program:cloudstrype]
user = cloudstrype
command = /usr/share/nginx/cloudstrype/venv/bin/uwsgi -c /usr/share/nginx/cloudstrype/deploy/uwsgi-cloudstrype.ini
directory = /usr/share/nginx/cloudstrype/
autostart = true
autorestart = true
stderr_logfile = /var/log/cloudstrype/uwsgi.stderr.log
stdout_logfile = /var/log/cloudstrype/uwsgi.stdout.log
stopsignal = INT

