# deploy/supervisord-array.ini
[program:array]
user = cloudstrype
command = /usr/share/nginx/cloudstrype/venv/bin/python3 manage.py array --bind=localhost --port=8765
directory = /usr/share/nginx/cloudstrype/web/cloudstrype/
autostart = true
autorestart = true
stderr_logfile = /var/log/cloudstrype/array.stderr.log
stdout_logfile = /var/log/cloudstrype/array.stdout.log
stopsignal = INT
