# cron/cron.d/fsck
*/5 * * * *	root	/usr/bin/env python3 /web/cloudstrype/manage.py fsck >> /var/log/cron.log 2>&1
#
